# docker/domserver/supervisor/nginx.conf
[program:nginx]
command=bash -c "sleep 10 && exec nginx -g 'daemon off;'"
numprocs=1
autostart=true
autorestart=true
